# build/image/mkosi.extra/etc/systemd/system/containerd.service.d/10-data-partition.conf
[Unit]
After=data.mount
Requires=data.mount
