# build/image-gpu/mkosi.extra/etc/systemd/system/manager.service.d/10-disk-mounter.conf
[Unit]
# Image-gpu variant: ensure cloud-optimised image and model disks are
# mounted before the manager starts. Without this, containers that
# reference disk:// or use ai.privasys.volume= would race with mount.
After=disk-mounter.service
Wants=disk-mounter.service
