# build/image/mkosi.extra/usr/lib/systemd/system/data.mount
[Unit]
Description=Data partition (/data)
After=luks-data.service
Requires=luks-data.service

[Mount]
What=/dev/mapper/data-crypt
Where=/data
Type=ext4
Options=defaults,noatime

[Install]
WantedBy=multi-user.target
