# build/image/mkosi.extra/usr/lib/systemd/system/home.mount
[Unit]
Description=Home directories (tmpfs)
DefaultDependencies=no
Before=local-fs.target

[Mount]
What=tmpfs
Where=/home
Type=tmpfs
Options=mode=0755,size=64M,nosuid,nodev

[Install]
WantedBy=local-fs.target
